# root/etc/services.d/rsyslogd/run
#!/usr/bin/execlineb -P

rsyslogd -f /etc/rsyslog.conf -n
